# systemd_units/oem-first-boot.service.d/00-mfr-util.conf
[Unit]
ConditionPathExists=/usr/bin/mfr_util
